# archive/shutter-island/tests/cron/GoProDownload
### BEGIN INIT INFO
# Provides:             telasocial
# Required-Start:
# Required-Stop:
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    TelaSocial Node App
### END INIT INFO

export PATH=$PATH:/usr/local/bin/node

case "$1" in
  start)
    /usr/local/bin/node /home/pi/ShutterIsland/repo/shutter-island/cloudberry/download.js
    ;;
  stop)
    exit 0
    ;;
  *)

  echo "Usage"
  exit 1
  ;;
esac
exit 0

