# systemd/cloudberry.service
[Unit]
Description=Cloudberry rig controller
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/cloudberry
ExecStart=%h/cloudberry/.venv/bin/cloudberry

[Install]
WantedBy=default.target
